# antlir/linux/metalos/workload-pre.target
[Unit]
Description=Host workload (Pre)
RefuseManualStart=yes
Requires=core-services.target
Requires=network-online.target
