# antlir/linux/metalos/workload-pre.target
[Unit]
Description=Host workload (Pre)
RefuseManualStart=yes
Before=workload.target
