# metalos/metald/metald.socket
[Socket]
ListenStream=/run/metalos/metald_socket
Accept=no
