# metalos/metald/metalos.metald.service
[Service]
ExecStart=/bin/echo "placeholder"
