# metalos/metald/metald.service
[Service]
Type=simple
ExecStart=/usr/lib/metalos/metald --systemd-socket
Restart=always
