# fs_image/bzl/tests/cheese-source.service
[Unit]
Description=Another Cheesy Test Service that is an image.source

[Service]
ExecStart=/bin/true
