# fs_image/bzl/tests/cheese-export.service
[Unit]
Description=Another Cheesy Test Service that is an explicit export_file

[Service]
ExecStart=/bin/true
