# metalos/os/tests/snapshot-test-control.service
[Unit]
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=btrfs subvolume snapshot /test-control-subvol-src /test-control-subvol
