# metalos/os/core-services.target
[Unit]
Description=Core services
After=core-services-pre.target
Requires=core-services-pre.target
AllowIsolate=yes
