# metalos/os/require-sshd-keygen.service
[Unit]
Requires=sshd-keygen.service
After=sshd-keygen.service
