# metalos/os/run-fs-control.mount
[Unit]
Description=Root Subvol Layout

[Mount]
Where=/run/fs/control
What=LABEL=/
Options=subvolid=5
