# metalos/os/workload-pre.target
[Unit]
Description=Host workload (Pre)
RefuseManualStart=yes
After=core-services.target
Requires=core-services.target network-online.target
