# metalos/os/workload-pre.target
[Unit]
Description=Host workload (Pre)
RefuseManualStart=yes
Requires=core-services.target network-online.target
