# antlir/vm/initrd/sysroot.mount
[Unit]
Description=Mount sysroot
After=local-fs.target
DefaultDependencies=no

[Mount]
What=LABEL=/
Where=/sysroot
Options=subvol=volume
