# antlir/vm/initrd/sysroot.mount
[Unit]
Description=vmtest root
Before=seedroot.service

[Mount]
Where=/sysroot
What=/dev/vda
Type=btrfs
Options=subvol=volume
