# antlir/vm/initrd/sysroot.mount
[Mount]
What=LABEL=/
Where=/sysroot
Options=subvol=volume
