# antlir/antlir2/antlir2_facts/tests/foo.service
[Unit]
Description=foo service

[Service]
Type=oneshot
ExecStart=true
