# antlir/bzl/tests/cheese-user.service
[Unit]
Description=Another Cheesy Test Service that is installed as a user service

[Service]
ExecStart=/bin/true
