# antlir/bzl/tests/cheese-requires.service
[Unit]
Description=Another Cheesy Test Service that is an implicit export_file

[Service]
ExecStart=/bin/true
