# antlir/bzl/tests/cheese-template@.service
[Unit]
Description=Another Cheesy Test Service that is a template unit. Instance name: %i

[Service]
ExecStart=/bin/true
