# antlir/antlir2/testing/tests/slow-unit.service
[Service]
Type=oneshot
ExecStart=sleep 2
