# antlir/compiler/test_images/fake-static.service
[Unit]
Description=Static fake service from bootable-systemd-os-with-buck-runnables

[Service]
ExecStart=/fake-service static
